# deploy/tagcloud@.service
# Template-unit для multi-instance деплоя tagcloud.
#
# Каждый instance — отдельный Node-процесс, слушающий свой порт.
# `%i` в имени unit'а (например, `tagcloud@3001`) подставляется как PORT.
#
# Установка (4 инстанса 3001–3004):
#   sudo cp deploy/tagcloud@.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now tagcloud@3001 tagcloud@3002 tagcloud@3003 tagcloud@3004
#
# Перед этим в Caddyfile добавить upstream-список — см. deploy/README.md
# раздел "Профиль multi-instance".
#
# ВАЖНО: одновременно с этим выключите single-instance:
#   sudo systemctl disable --now tagcloud

[Unit]
Description=tagcloud (instance %i)
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=tagcloud
Group=tagcloud
WorkingDirectory=/opt/tagcloud
EnvironmentFile=/etc/tagcloud/tagcloud.env
# Перебиваем PORT из tagcloud.env значением из имени unit'а (после @).
Environment=PORT=%i
ExecStart=/usr/bin/node /opt/tagcloud/build/index.js

Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# graceful shutdown (см. tagcloud.service)
TimeoutStopSec=30
KillSignal=SIGTERM

# Hardening (идентичен single-instance — централизовать через drop-in
# нельзя, т.к. для template-unit это требует отдельного override-файла).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/tagcloud
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
